# reference/pi-listener-python/openvend-listener.service
[Unit]
Description=Open MDB MDB vend listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openvend
WorkingDirectory=/home/openvend/openvend/reference/pi-listener-python
EnvironmentFile=/home/openvend/openvend/reference/pi-listener-python/.env
ExecStart=/home/openvend/openvend/reference/pi-listener-python/.venv/bin/python main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
